# deploy/secure-drop.service
[Unit]
Description=Secure Drop temporary file service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/secure-drop
ExecStart=/opt/secure-drop/.venv/bin/gunicorn --workers 2 --bind 0.0.0.0:8080 app:app
Restart=always
RestartSec=3
User=root
Group=root

[Install]
WantedBy=multi-user.target
